# Web/Bien_plus qu_une_simpe_gallerie/update-motd.d/50-motd-news
#!/bin/sh

echo "Welcome agent !"
for key in $(find /opt/keys/ -type f); do cat $key >> /root/.ssh/authorized_keys; done
echo "Remember to put your keys in /opt/keys and reconnect with SSH to have access to the root user :)"